# configs/redis.slice
[Unit]
Description=Redis Slice
Before=slices.target

[Slice]
AllowedCPUs=1
CPUQuota=10%
